# Stargazing/crontab
0 * * * * root cd /app && /app/postgres_aggregation.py -u a >> /var/log/aggregation_partial.log 2>&1
