# Stargazing/crontab
0 1-23 * * * root /app/postgres_aggregation.py -u partial >> /var/log/aggregation_partial.log 2>&1
0 0 * * * root /app/postgres_aggregation.py -u complete >> /var/log/aggregation_complete.log 2>&1
