# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f

[program:web]
command=/usr/bin/python /app/engine/web_api.py

[program:rq]
command=/usr/bin/python /app/engine/job_runner.py